# power/99-power-profile.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/power-profile-policy.sh"
