# power/power-profile-boot.service
[Unit]
Description=Set power profile at boot based on AC state
After=power-profiles-daemon.service
Requires=power-profiles-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-profile-policy.sh

[Install]
WantedBy=multi-user.target
